# MiniEPG.crontab
# Install this file into crontab:
#  $ crontab $HOME/MiniEPG/MiniEPG.crontab
#
40 6,18 * * * bash $HOME/MiniEPG/MiniEPG.sh
