# MiniEPG.crontab
# Install this file into crontab:
#  $ crontab $HOME/MiniEPG/MiniEPG.crontab
#
4 8,20 * * * bash $HOME/MiniEPG/MiniEPG_cron.sh

